# deps/iperf3-server.service
[Unit]
Description=Iperf3 Server

[Service]
Type=simple
# Requires the "expect" package for "unbuffer"
ExecStart=/usr/bin/unbuffer /usr/bin/iperf3 -s

[Install]
WantedBy=multi-user.target
