# deps/phc-to-phc-sync.service
[Unit]
Description=Synchronize one PTP Hardware Clock to another
Requires=ptp4l.service
After=ptp4l.service

[Service]
Type=simple
ExecStart=/usr/sbin/phc2sys -s /dev/ptp1 -c /dev/ptp0 -O 0

[Install]
WantedBy=multi-user.target
