# deps/phc-to-phc-sync.service
[Unit]
Description=Synchronize one PTP Hardware Clock to another

[Service]
Type=simple
ExecStart=/usr/sbin/ptp4l -i #IFACE# -f #CFG# --BMCA noop --logSyncInterval -3 --network_transport L2 --delay_mechanism P2P --tx_timestamp_timeout 10 --uds_address /var/run/ptp4l-loopback

[Install]
WantedBy=multi-user.target
