# deps/ptp4l.service
[Unit]
Description=Precision Time Protocol (PTP) service

[Service]
Type=simple
ExecStartPre=/sbin/hwclock --hctosys
ExecStartPre=/usr/sbin/phc_ctl /dev/ptp1 set
ExecStartPre=/usr/sbin/phc_ctl /dev/ptp1 freq 0
ExecStart=/usr/sbin/ptp4l -f /etc/linuxptp/ptp4l.conf -i swp1

[Install]
WantedBy=multi-user.target
